# orchestra.service
[Unit]
Description=Orchestra FastAPI Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/orchestra/backend
Environment="PATH=/root/orchestra/venv/bin"
EnvironmentFile=/root/orchestra/backend/.env
ExecStart=/root/orchestra/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
